# pi/crontabs/user_pi
# Crontab aus https://github.com/ac-caterva/webserver/pi/crontabs/user_pi
#
# Date: 2021-01-10
#
# m h  dom mon dow   command
#
# Copy Caterva log file invoiceLog.csv every minute
* * * * * /var/caterva/scripts/copy_log.sh 
# Update Pi and Caterva if requested
* * * * * /var/caterva/scripts/update_Pi_and_Caterva.sh
